# helper/supervisord.conf
[unix_http_server]
file = /var/run/supervisor.sock
username = root
password = toor

[inet_http_server]
port = 9001
username = root
password = toor

[supervisord]
;user=supervisor
loglevel = info
logfile = /var/log/supervisord.log
pidfile = /var/run/supervisord.pid
; FIX : need to delete if apache and mysql are cool 
childlogdir = /var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = http://localhost:9001
username = root
password = toor
history_file=~/.sc_history

;[include]
;files = /etc/supervisor/conf.d/*.conf

[program:creep]
command=/opt/creepMiner/bin/creepMiner
numprocs=1
autostart=true
autorestart=true